# system/SITE_NAME_PLACEHOLDER.socket
[Unit]
Description=SITE_NAME_PALCEHOLDER socket

[Socket]
ListenStream=/run/SITE_NAME_PALCEHOLDER.sock

[Install]
WantedBy=sockets.target