# .system/SITE_NAME_PLACEHOLDER.socket
[Unit]
Description=SITE_NAME_PLACEHOLDER socket

[Socket]
ListenStream=/run/SITE_NAME_PLACEHOLDER.sock

[Install]
WantedBy=sockets.target